# jenkins/cron_job/crontab
# crontab -e
*/2 * * * * root cd /usr/local/bin/ && python3 check_workspace.py