# jenkins/logging/cron_job/crontab
# crontab -e
2-59/2 * * * * cd /usr/local/bin/ && command python3 check_workspace.py